# alpine/etc/s6/x11vnc/finish
#!/bin/sh

# Terminate all services.
s6-svscanctl -t /etc/s6
